# scripts/com.cuadram.claude-config-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cuadram.claude-config-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/cuadram/proyectos/backup_claude_config.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/cuadram/Library/Logs/claude-config-backup.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/cuadram/Library/Logs/claude-config-backup.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
